# templates/etc/systemd/system/sentry-web.service.j2
# {{ ansible_managed }}
[Unit]
Description=Sentry Main Service
After=network.target
Requires=sentry-worker.service
Requires=sentry-cron.service

[Service]
Type=simple
User={{ sentry__user }}
Group={{ sentry__group }}
WorkingDirectory={{ sentry__home }}
Environment=SENTRY_CONF=/etc/sentry
ExecStart={{ sentry__virtualenv }}/bin/sentry run web

[Install]
WantedBy=multi-user.target
